# twix.cron
* * * * * ./home/josiah/twixpenser/twix.rb
